# root/etc/services.d/pulseaudio/run
#!/usr/bin/with-contenv bash
exec \
	s6-setuidgid abc pulseaudio --daemonize=0 > /dev/null 2>&1
