# root/etc/services.d/kodi/run
#!/usr/bin/with-contenv bash

	s6-setuidgid abc \
	/usr/bin/kodi --headless \
	> /config/.kodi/kodi.screen_output.log 2>&1
