# root/etc/services.d/kodi/run
#!/usr/bin/with-contenv bash
s6-setuidgid abc /app/kodi/lib/kodi/kodi-x11 --headless
